# installs/docker/OSP-Edge/docker-files.d/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/nginx/sbin/nginx -g "daemon off;"