# installs/docker/OSP-Core/docker-files.d/supervisord-celery.conf
[supervisord]
nodaemon=true

[program:osp-celery]
directory=/opt/osp
user=www-data
group=www-data
command=/usr/local/bin/celery -A app.celery worker --time-limit=3600 --concurrency=4